# Nayel Application/EV-Backend-0.2/EV-Backend-0.2/.platform/files/supervisor.ini
# The following given example is used with Horizon, but you can
# make a copy of it or add your own and configure it.

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/app/current/artisan queue:work --queue=accountVerification,forgetPassword,appLoginCredentials --sleep=3 --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/worker.log
stopwaitsecs=3600
